# root/etc/services.d/authelia/run
#!/bin/sh
# Authelia s6 service

echo "[INFO] Starting authelia..."

# Standard start command for authelia
exec s6-setuidgid bsd /app/authelia --config /config/authelia.yml
